# build/ws_supervisor.conf
[program:ws-worker]
directory=/home/badamczyk/webtetrado
command=uvicorn webTetrado.asgi:application --reload --port 8081
autostart=true
autorestart=true
user=root
numprocs=1
startretries=50
stopsignal=INT
killasgroup=true
stopasgroup=true
stdout_logfile=/home/badamczyk/webtetrado/supervisor/logs.conf
redirect_stderr=true